# WLS.service
[Unit]
Description=Wireless Light Sensor (WLS) — Android ALS to Monitor Brightness
Documentation=https://github.com/prasanthc41m/WLS
After=graphical-session.target network-online.target
Wants=network-online.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/WLS.py
Restart=on-failure
RestartSec=5
# Wayland environment
Environment=WAYLAND_DISPLAY=wayland-0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=QT_QPA_PLATFORM=wayland

[Install]
WantedBy=graphical-session.target
